# cloudtunes-server/production/supervisor.ini
[group:cloudtunes-server]
programs=cloudtunes-8000,cloudtunes-8001


[program:cloudtunes-8000]
command=cloudtunes-server --port=8000
directory=/home/cloudtunes/cloudtunes-server
user=cloudtunes
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/cloudtunes-server.log
loglevel=debug


[program:cloudtunes-8001]
command=cloudtunes-server --port=8001
directory=/home/cloudtunes/cloudtunes-server
user=cloudtunes
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/cloudtunes-server.log
loglevel=debug


[program:cloudtunes-worker]
command=cloudtunes-worker worker -c 50
directory=/home/cloudtunes/cloudtunes-server
user=cloudtunes
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/cloudtunes-worker.log
loglevel=debug
